# ufd-exporter.timer
[Unit]
Description=ufd Exporter Timer

[Timer]
OnCalendar=Weekly
Persistent=true

[Install]
WantedBy=timers.target